# Certbot/certbot.service
[Unit]
Description=Run certbot 30 minutes before expiring

[Service]
Type=oneshot
ExecStart=/root/renew-cert.sh
User=root


